# kpi-web.service
[Unit]
Description=Lotus KPI Web Application
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/kpi-web
Environment=SECRET_KEY=change-this-to-a-long-random-string
Environment=KPI_ADMIN_USER=admin
Environment=KPI_ADMIN_PASS=admin
ExecStart=/opt/kpi-web/venv/bin/gunicorn -w 1 -b 0.0.0.0:20000 app:app --timeout 120
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
